# mutt/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c '/usr/bin/mbsync -Va -c ~/.config/mutt/mbsyncrc'
